# Library/LaunchAgents/botanicus.cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>botanicus.cron</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
    	<string>/Users/botanicus/.gem/ruby/2.4.0/bin:/Users/botanicus/.rubies/ruby-2.4.0/bin:/Users/botanicus/.scripts:/bin:/usr/local/bin:/usr/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/botanicus/.bitbar</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/botanicus/.bitbar/bin/cron</string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <true />

    <!-- I don't have permissions to write into /var/log. -->
    <key>StandardOutPath</key>
    <string>/tmp/botanicus.cron.stdout</string>
    <key>StandardErrorPath</key>
    <string>/tmp/botanicus.cron.stderr</string>
  </dict>
</plist>
