# Library/LaunchAgents/botanicus.commit_and_push_projects.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>botanicus.commit_and_push_projects</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/zsh</string>
      <string>-c</string>
      <string>launchd.commit_and_push_projects.sh</string>
    </array>

    <!-- Run at 18:20 every day. -->
    <key>StartCalendarInterval</key>
    <dict>
    	<key>Hour</key>
    	<integer>18</integer>
    	<key>Minute</key>
    	<integer>20</integer>
    </dict>

    <!-- I don't have permissions to write into /var/log. -->
    <key>StandardOutPath</key>
    <string>/tmp/botanicus.commit_and_push_projects.stdout</string>
    <key>StandardErrorPath</key>
    <string>/tmp/botanicus.commit_and_push_projects.stderr</string>
  </dict>
</plist>
